# lib/golden/agent/claude-agent-morty.service
# /etc/systemd/system/claude-agent-morty.service
# Managed by bubble-vps-platform (SPEC-007 + SPEC-009). Do not edit by hand.
[Unit]
Description=Claude Agent — morty (bubble-internal)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/home/claude/agents/morty

# CRITICAL: claude spawns plugin processes via `bun run`. systemd's default
# PATH does NOT include /home/claude/.bun/bin where bun lives, so without
# this override the Telegram plugin (and any other bun-based plugin) fails
# silently — claude looks for `bun` on PATH, doesn't find it, and never
# spawns the plugin process. Symptom: service is "active", but no plugin
# child processes, no Telegram polling, messages pile up in pending_update_count.
# Fix discovered 2026-05-09 after Joris reported "Morty doesn't receive texts".
# Same gotcha as documented in projects/hetzner-migration/CLAUDE_AGENT_RECIPE.md
# §"The four gotchas" #1 — applies to systemd just as it does to MCP child
# spawning. settings.json's env.PATH is for claude's INTERNAL use; child-process
# environment comes from the parent service's environment.
Environment=PATH=/home/claude/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/bin

# Per-concierge Telegram channel-state dir (SPEC-021 multi-agent finding).
# The Telegram MCP plugin (server.ts) keeps ALL of its per-agent runtime state
# — access.json, approved/, inbox/, .env AND bot.pid — under the directory
# named by TELEGRAM_STATE_DIR, defaulting to ~/.claude/channels/telegram when
# the var is UNSET. On a multi-concierge box (morty + claudette on joris-cx33)
# leaving it unset makes BOTH pollers default to the same bare telegram/ dir →
# they clobber each other's bot.pid and fight over the single getUpdates
# long-poll slot (Telegram 409 Conflict). So each unit exports its OWN dir.
# telegram_state_dir is derived (in _systemd_unit.py) from the SAME
# lib.host_helpers.telegram_channel_dir single source that _telegram_plugin.py
# uses to CREATE the dir — so the exported var and the created dir can never
# diverge. morty (the original/default agent) keeps the BARE telegram/ dir so
# its existing live channel state isn't orphaned; others get telegram-<name>/.
Environment=TELEGRAM_STATE_DIR=/home/claude/.claude/channels/telegram

# Decrypt secrets into tmpfs at startup. The `+` prefix makes ExecStartPre
# run as root (overriding User=claude) so it can read the age private key
# at /etc/age/key.txt. The decrypted file is then chowned to claude so
# the service process can read it.
#
# runtime_env_dir is PER-CONCIERGE (SPEC-001 v1.2 multi-concierge): the primary
# concierge keeps the historical /run/claude-agent (so morty's live unit does
# not churn), additional concierges decrypt into their own /run/claude-agent-<name>
# so two concierges on one box never clobber each other's tmpfs env file.
ExecStartPre=+/bin/mkdir -p /run/claude-agent
ExecStartPre=+/bin/chown claude:claude /run/claude-agent
ExecStartPre=+/bin/chmod 0750 /run/claude-agent
ExecStartPre=+/bin/sh -c 'SOPS_AGE_KEY_FILE=/etc/age/key.txt /usr/local/bin/sops --decrypt --output /run/claude-agent/env.tmp /etc/bubble/secrets.sops.env'
ExecStartPre=+/bin/mv /run/claude-agent/env.tmp /run/claude-agent/env
ExecStartPre=+/bin/chmod 0400 /run/claude-agent/env
ExecStartPre=+/bin/chown claude:claude /run/claude-agent/env

# `-` prefix on EnvironmentFile makes it optional at unit-load time. The
# file is created by the ExecStartPre chain above. Without the `-`, systemd
# tries to load the file at the very first activation step and fails with
# "Failed to load environment files: No such file or directory" before
# ExecStartPre ever runs (observed on Ubuntu 24.04 + systemd 255).
EnvironmentFile=-/run/claude-agent/env
# claude detects no-TTY and switches to --print mode (which then errors on
# missing prompt). The interactive REPL needs a pty. `script -qfc CMD /dev/null`
# allocates a pty, runs CMD inside it, and discards typescript output. This
# matches what tmux did in the legacy setup, but as a foreground process
# systemd can supervise (Type=simple).
#
# Permission mode: --dangerously-skip-permissions (canonical alias for
# bypassPermissions per Anthropic docs §"Skip all checks"). This mode CANNOT
# be activated from settings.json defaultMode alone — it MUST be passed as a
# CLI flag at session start.
#
# Why this is acceptable here: joris-cx33 is an isolated single-tenant VM
# (Hetzner CX33). Per the doc's own guidance: "Only use this mode in isolated
# environments like containers, VMs, or dev containers". auto mode is
# unreachable due to TTY-bound opt-in; acceptEdits still prompts on every
# Bash; this flag is the only way to fully eliminate prompts in headless ops.
#
# Wrapped in /bin/sh -c "..." because systemd's ExecStart parser splits on
# whitespace inside single-quoted args (verified 2026-05-09: without the sh
# wrapper, systemd unquoted the script -qfc argument and passed each token
# as a separate argv to script, breaking the launch).
ExecStart=/bin/sh -c '/usr/bin/script -qfc "/usr/bin/claude --dangerously-skip-permissions --channels plugin:telegram@claude-plugins-official" /dev/null'

Restart=on-failure
RestartSec=10
LimitNOFILE=65536

# Decrypted env file is removed when service stops — defense in depth even
# though /run is tmpfs.
ExecStopPost=+/bin/rm -f /run/claude-agent/env

[Install]
WantedBy=multi-user.target
